# com.trinity.vr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trinity.vr</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/tybrown/Desktop/Trinity-System/vr_server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tybrown/Desktop/Trinity-System</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tybrown/Desktop/Trinity-System/logs/vr_server_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tybrown/Desktop/Trinity-System/logs/vr_server_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>Nice</key>
    <integer>-10</integer>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
